# sdk/config/arm/raspberrypi0/overlay/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Mount filesystems
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mkdir -p /dev/shm
::sysinit:/bin/mount -a

# Initialize /var
::sysinit:/bin/mkdir -p /var/log
::sysinit:/bin/mkdir -p /var/lib/mediabox
::sysinit:/bin/ln -sf /data /var/lib/mediabox/store
::sysinit:/bin/mkdir -p /var/lib/mediabox/store/Audio
::sysinit:/bin/mkdir -p /var/lib/mediabox/store/Video
::sysinit:/bin/mkdir -p /var/lib/dropbear
::sysinit:/bin/ln -s /tmp /var/run

# Network setup
::sysinit:/bin/hostname -F /boot/hostname
::sysinit:/sbin/ifup -a
::sysinit:/usr/sbin/wpa_supplicant -B -Dnl80211 -iwlan0 -c /etc/wpa_supplicant.conf
::sysinit:/sbin/udhcpc -b -iwlan0
::sysinit:/usr/sbin/ntpd -q -p ptbtime1.ptb.de
::sysinit:/usr/sbin/dropbear -R -B

# Rotate logs and start mediabox
::sysinit:/bin/mv -f /var/log/mediabox.log /var/log/mediabox.log.1 2>&1 > /dev/null
::sysinit:/bin/mv -f /var/log/avmount-mediabox.log /var/log/avmount-mediabox.log.1 2>&1 >/dev/null
::sysinit:/sbin/start-stop-daemon -S -q -x /usr/bin/mediabox -p /tmp/mediabox.pid -b -m -- --avbox:logfile /var/log/mediabox.log --no-mediatomb

# now run any rc scripts
::sysinit:/etc/init.d/rcS

# Put a getty on the serial port
tty::respawn:/sbin/getty -L console 0 vt100 # Default console

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcK
::shutdown:/sbin/start-stop-daemon -K -p /tmp/mediabox.pid
::shutdown:/bin/umount -a -r
